# contrib/backer
#!/sbin/openrc-run

name="$SVCNAME"
command="/usr/local/sbin/$SVCNAME"
supervisor="supervise-daemon"
pidfile="/run/$SVCNAME/$SVCNAME.pid"

depend() {
        need net
}

start_pre() {
        mkdir -p /var/log/$SVCNAME
        mkdir -p /run/$SVCNAME
}

# vim: set ft=sh noet ai ts=4 sw=4 sts=4: